# bin/grid-monitor.service
[Unit]
Description=Power Grid Monitor
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/srv/tasks/bin/grid-monitor.sh
User=toni
Group=toni

[Install]
WantedBy=multi-user.target
